# infrastructure/systemd/fullstack-app.service
[Unit]
Description=Full-Stack Application
Documentation=https://github.com/yourusername/fullstack-boilerplate
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=appuser
Group=appuser
WorkingDirectory=/opt/fullstack-app
RemainAfterExit=yes

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/opt/fullstack-app/.env

# Start command
ExecStart=/opt/fullstack-app/infrastructure/scripts/deploy-prod.sh

# Stop command (clean up containers)
ExecStop=/bin/bash -c 'podman rm -f fullstack-postgres fullstack-backend fullstack-frontend fullstack-nginx 2>/dev/null || true'

# Restart policy
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=60

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/fullstack-app

[Install]
WantedBy=multi-user.target
